# pve-network-rs/services/pve-network-rs.service
[Unit]
Description=Proxmox VE Network Management Service (Rust)
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/bin/pve-network-rs-daemon
Restart=always
RestartSec=5
User=root
Group=root

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/network /etc/pve /run /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target